# new-device-notify/new-device-notify.timer
[Unit]
Description=Scan for new devices in the network on the schedule

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
