# twitch-live-notify/twitch-live-notify.service
[Unit]
Description=Twitch Channel Live notifier

[Service]
Type=simple
ExecStart=php %h/.local/bin/twitch-live-notify
