# !abandoned/new-device-notify/new-device-notify.timer
[Unit]
Description=Regularly check for new devices in the network

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
