# showRSS/showRSS.timer
[Unit]
Description=Run showRSS checker regularly

[Timer]
OnCalendar=*:0/10
RandomizedDelaySec=60s

[Install]
WantedBy=timers.target
