# showRSS/showRSS.service
[Unit]
Description=ShowRSS checker and notifier

[Service]
Type=simple
ExecStart=%h/.local/bin/showRSS
