# fmovies-watchlist/fmovies-watchlist.service
[Unit]
Description=FMovies watchlist notifier

[Service]
Type=simple
ExecStart=php %h/.local/bin/fmovies-watchlist
