# fmovies-watchlist/fmovies-watchlist.timer
[Unit]
Description=Run FMovies watchlist notifier regularly

[Timer]
OnCalendar=*:0/10
RandomizedDelaySec=60s

[Install]
WantedBy=default.target
