# dst-reward-links/dst-reward-links.timer
[Unit]
Description=Run DST Reward Links notifier

[Timer]
OnCalendar=00/1:00
RandomizedDelaySec=60s

[Install]
WantedBy=default.target
